# guesthandler.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='system/guesthandler' type='service' version='0'>
    <create_default_instance enabled='true'/>
    <single_instance/>
    <dependency name='vmadmd' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/smartdc/vmadmd:default'/>
    </dependency>
    <exec_method name='start' type='method' exec='%{config/method_prefix}/svc-guesthandler %m' timeout_seconds='0'>
      <method_context>
        <method_environment>
          <envvar name='DELAY' value='30'/>
        </method_environment>
      </method_context>
    </exec_method>
    <exec_method name='stop' type='method' exec='%{config/method_prefix}/svc-guesthandler %m' timeout_seconds='0'>
      <method_context>
        <method_environment>
          <envvar name='GRACE_PERIOD' value='300'/>
          <envvar name='WAIT_INTERVAL' value='10'/>
        </method_environment>
      </method_context>
    </exec_method>
    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='transient'/>
    </property_group>
    <property_group name='config' type='eu.gerczei,guesthandler_configuration'>
      <propval name='method_prefix' type='astring' value='/opt/custom/bin'/>
    </property_group>
    <stability value='Unstable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Ordered start-up and shutdown for zones</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
